# templates/mbsync.timer.j2
[Unit]
Description=Sync mail periodically

[Timer]
OnCalendar={{ mail_check_schedule }}
Persistent=true

[Install]
WantedBy=default.target
