# templates/mbsync.timer.j2
[Unit]
Description=Sync mail periodically

[Timer]
OnUnitInactiveSec={{ mail_interval|default('4min') }}
Persistent=true

[Install]
WantedBy=default.target
