# templates/mbsync.service.j2
[Unit]
Description=Sync mail

[Service]
ExecStart=/usr/bin/mbsync -aV
Type=oneshot
