# templates/mbsync.service.j2
[Unit]
Description=Sync mail

[Service]
ExecStart=/bin/sh -c "mbsync -aV && exec notmuch new"
Type=simple
Restart=on-failure
RestartSec={{ mail_restart_interval }}
RuntimeMaxSec={{ mail_timeout }}
