# templates/notmuch-compact.timer.j2
[Unit]
Description=Compact notmuch periodically

[Timer]
OnCalendar={{ mail_compact_schedule }}
Persistent=true

[Install]
WantedBy=default.target
