# templates/notmuch-compact.service.j2
[Unit]
Description=Compact notmuch database

[Service]
ExecStart=/usr/bin/notmuch compact
Type=simple
TimeoutStartSec=1h
